# system/send-ip.service
[Unit]
Description=Updating IP address of a server.

[Service]
Type=simple
ExecStart=/usr/local/bin/send-ip.sh
